# backend/crontab
# m h  dom mon dow   command
*/5 * * * * /home/ubuntu/betterspacecal/backend/release.sh
42 * * * * /home/ubuntu/betterspacecal/backend/tweet.py
0 0 * * * /home/ubuntu/betterspacecal/backend/daily.sh
